# net/softether/etc/systemd/system/softether-vpnserver.service.d/override.conf
[Service]
ExecStartPost=/opt/vpnserver/post_start.sh
